# systemd/granturismo.service
[Unit]
Description=GranTurismo Online (web game + leaderboard API)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/http/granturismo
ExecStart=/usr/bin/node server/server.js
Environment=PORT=9000
Environment=NODE_ENV=production
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
